# protected/bin/linux_bash/supervisor.ini

[program:node5]
command=node /data/webapps/spider.duowan.com/protected/bin/checkProxyPool.js
process_name=checkProxyPool
directory=/data/webapps/spider.duowan.com/protected/bin/
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/tmp/checkProxyPool.log

[program:node6]
command=node /data/webapps/spider.duowan.com/protected/bin/crawl.js
process_name=spider_crawl
directory=/data/webapps/spider.duowan.com/protected/
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/tmp/spider_crawl.log
